# internal/assets/files/systemd/olcrtc-wb-session.service
[Unit]
Description=olcRTC WB Stream Playwright session
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=olcrtc-wb
Group=olcrtc-wb
UMask=0077
RuntimeDirectory=olcrtc-wb
RuntimeDirectoryMode=0750
RuntimeDirectoryPreserve=yes
StateDirectory=olcrtc-wb
StateDirectoryMode=0700
ExecStart=/bin/bash /usr/lib/olcrtc-panel/wb/run-session.sh
KillMode=control-group
TimeoutStopSec=10
MemoryMax=1200M
TasksMax=256
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths=/var/lib/olcrtc-wb /run/olcrtc-wb
StandardOutput=journal
StandardError=journal
